# dashboard_service/supervisord.conf
[supervisord]
nodaemon=true

[program:dashboard_service]
command=/app/start_server.sh

[program:update_dashboard_worker]
command=celery -A dashboard_service worker -l info --pool=solo -Q update_dashboard_topic

[program:update_dashboard_beat]
command=celery -A dashboard_service beat -l info